# packaging/debian/tortunnel-helper.service
[Unit]
Description=TorTunnel privileged helper
Documentation=man:tortunnel-helper(8)
After=network-online.target

[Service]
Type=simple
Environment=TORTUNNEL_COUNTRY=DE
EnvironmentFile=-/etc/tortunnel/helper.env
ExecStart=/usr/libexec/tortunnel/tortunnel-helper up --country ${TORTUNNEL_COUNTRY}
# Fail-closed teardown: remove the nftables table and TUN device on stop/crash.
ExecStopPost=/usr/libexec/tortunnel/tortunnel-helper down
Restart=on-failure
RestartSec=2s
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/tortunnel /var/lib/tortunnel

[Install]
WantedBy=multi-user.target
